# rc/zsh/profile.d/gvm.zsh
if [ -r $HOME/.gvm/scripts/gvm ]; then
    source $HOME/.gvm/scripts/gvm
fi
