# rc/zsh/profile.d/plenv.zsh
_glob_add_path_env PATH "${HOME}/.plenv/bin"
