# modules/zsh/dotfiles/zsh/profile.d/go.zsh
export GOPATH=$HOME/projects/go
export GOBIN=$GOPATH/bin
_glob_add_path_env PATH $GOBIN
