# modules/zsh/dotfiles/zsh/profile.d/rbenv.zsh
_glob_add_path_env PATH "${HOME}/.rbenv/bin"
