# etc/zsh/profile.d/00_functions
_glob_source "$ZDOTDIR/function.d/*.zsh"
