# etc/zsh/profile.d/android
SDK_ROOT=$HOME/opt/pkgs/android-sdk-mac_x86/
export SDK_ROOT
