# etc/zsh/profile.d/courier-imap.sh
#
# courier-imap の設定
#

# PATH の設定
add_env PATH "/usr/lib/courier-imap/*bin"
