# etc/zsh/profile.d/cvs.sh
#
# CVS の設定
#

CVS_RSH=ssh
CVSROOT=:ext:cvs.bz2.jp:/home/masato/data/cvs/repository
export CVS_RSH CVSROOT
