# etc/zsh/profile.d/homebrew
# check homebrew installed
type brew > /dev/null 2>&1
if [ $? -ne 0 ]; then
  return
fi

# coreutils
if [ -d $(brew --prefix coreutils)/libexec/gnubin ]; then
  add_env PATH $(brew --prefix coreutils)/libexec/gnubin
fi
if [ -d $(brew --prefix coreutils)/libexec/gnuman ]; then
  add_env MANPATH $(brew --prefix coreutils)/libexec/gnuman
fi

# findutils
if [ -d $(brew --cellar findutils) ]; then
  alias find=gfind
  alias xargs=gxargs
fi

# gnu-sed
if [ -d $(brew --cellar gnu-sed) ]; then
  alias sed=gsed
fi

# z
if [ -f $(brew --prefix`/etc/profile.d/z.sh) ]; then
function precmd () {
  z --add "$(pwd -P)"
}
