# etc/zsh/profile.d/homebrew
_command_exists brew || return

# coreutils
add_env PATH $(brew --prefix coreutils)/libexec/gnubin
add_env MANPATH $(brew --prefix coreutils)/libexec/gnuman

# gnu-sed
add_env PATH $(brew --prefix gnu-sed)/libexec/gnubin
add_env MANPATH $(brew --prefix gnu-sed)/libexec/gnuman
