# etc/zsh/profile.d/homebrew
# check homebrew installed
type brew > /dev/null 2>&1
if [ $? -ne 0 ]; then
  return
fi

# coreutils
add_env PATH $(brew --prefix coreutils)/libexec/gnubin
add_env MANPATH $(brew --prefix coreutils)/libexec/gnuman

# gnu-sed
add_env PATH $(brew --prefix gnu-sed)/libexec/gnubin
add_env MANPATH $(brew --prefix gnu-sed)/libexec/gnuman
