# etc/zsh/profile.d/keychain
if ! which keychain > /dev/null 2>&1; then
  return
fi

keychain id_rsa id_dsa

if [ -f $HOME/.keychain/$HOST-sh ]; then
  source $HOME/.keychain/$HOST-sh
fi
