# etc/zsh/profile.d/mh
#
# mh の設定
#

# PATH の設定
add_env PATH "/usr/local/bin/mh" "/usr/bin/mh"
