# etc/zsh/profile.d/nodebrew
add_env PATH $HOME/.nodebrew/current/bin
