# etc/zsh/profile.d/node.js.sh
NODE_PATH=$HOME/opt/homebrew/lib/node
export NODE_PATH
