# modules/zsh/dot.zsh/profile.d/ndenv.zsh
_glob_add_path_env PATH "${HOME}/.ndenv/bin"
